# vendors/OpenGear/romfs/libexec/udev/rules.d/60-mtu.rules
ACTION=="add|change|remove", SUBSYSTEM=="net", RUN+="/etc/scripts/udev-mtu %k"
