# vendors/OpenGear/romfs/libexec/udev/rules.d/80-usbmodem.rules
SUBSYSTEM=="tty", SUBSYSTEMS=="usb", DRIVERS=="cdc_acm", RUN+="/etc/scripts/udev-usbmodem"

