# vendors/OpenGear/IM72xx/romfs/libexec/udev/rules.d/70-net.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="eth2", RUN+="/etc/scripts/udev-net-eth2"
