# vendors/Digi/EX50/76-power_supply_update_led.rules
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/etc/power-led.sh"
